# systemd/cracked-checker.service
[Unit]
Description=Cracked Alert System - Watchdog Checker
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/cracked_checker.sh
Restart=always
RestartSec=5

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cracked-checker

[Install]
WantedBy=multi-user.target